# FileSort v-Current-W-L/Configurator/Configurator/bin/Debug/net6.0/Scripts/sortingsystem.service
[Unit]
Description= Service for Sortingsystem, published 29.05.2023
[Service]
User=root
Type=notify
WorkingDirectory=/opt/SortingSystem/
ExecStart='/opt/SortingSystem/SortingSytem V-26-05-23'
[Install]
WantedBy=multi-user.target

